# desktop/why-cycle.service
[Unit]
Description="Changes the user's themes to match the current time of day."

[Service]
ExecStart=/usr/bin/why-cycle -r
Type=oneshot

